# templates/systemd/promtail.service.j2
{#
SPDX-FileCopyrightText: 2024 Slavi Pantaleev
SPDX-FileCopyrightText: 2024 Tiz
SPDX-FileCopyrightText: 2026 Suguru Hirahara

SPDX-License-Identifier: AGPL-3.0-or-later
#}

[Unit]
Description=Promtail log agent ({{ promtail_identifier }})
{% for service in promtail_systemd_required_services_list %}
Requires={{ service }}
After={{ service }}
{% endfor %}
{% for service in promtail_systemd_wanted_services_list %}
Wants={{ service }}
{% endfor %}
DefaultDependencies=no

[Service]
Type=simple
ExecStartPre={{ devture_systemd_docker_base_host_command_sh }} -c 'i=0; while [ "$i" -lt 10 ]; do {{ devture_systemd_docker_base_host_command_docker }} rm -f {{ promtail_identifier }} >/dev/null 2>&1 || true; if ! {{ devture_systemd_docker_base_host_command_docker }} inspect --type=container {{ promtail_identifier }} >/dev/null 2>&1; then exit 0; fi; i=$((i + 1)); sleep 1; done; echo "Failed to remove stale container {{ promtail_identifier }}" >&2; exit 1'

{#
    When running with a different user than promtail_uid/promtail_gid, even if that user is root,
    we need to grant CAP_DAC_READ_SEARCH so that Promtail can read its own config.
#}
ExecStartPre={{ devture_systemd_docker_base_host_command_docker }} create \
      --rm \
      --name={{ promtail_identifier }} \
      --log-driver=none \
      {% if promtail_container_run_as_non_root_user_enabled %}
      --user={{ promtail_uid }}:{{ promtail_gid }} \
      {% endif %}
      --cap-drop=ALL \
      {% for capability in promtail_container_capabilities_to_add %}
      --cap-add={{ capability }} \
      {% endfor %}
      --read-only \
      --network={{ promtail_container_network }} \
      {% if promtail_container_http_host_bind_port %}
      -p {{ promtail_container_http_host_bind_port }}:{{ promtail_config_http_listen_port }} \
      {% endif %}
      {% if promtail_container_grpc_host_bind_port %}
      -p {{ promtail_container_grpc_host_bind_port }}:{{ promtail_config_grpc_listen_port }} \
      {% endif %}
      --env-file={{ promtail_base_path }}/env \
      --label-file={{ promtail_base_path }}/labels \
      --mount type=bind,src={{ promtail_config_path }},dst=/etc/promtail,ro \
      --mount type=bind,src={{ promtail_data_path }},dst=/data \
      {% for mount in promtail_container_additional_volumes %}
      --mount {{ mount }} \
      {% endfor %}
      {% for arg in promtail_container_extra_arguments %}
      {{ arg }} \
      {% endfor %}
      {{ promtail_container_image }}

{% for network in promtail_container_additional_networks %}
ExecStartPre={{ devture_systemd_docker_base_host_command_docker }} network connect {{ network }} {{ promtail_identifier }}
{% endfor %}

ExecStart={{ devture_systemd_docker_base_host_command_docker }} start --attach {{ promtail_identifier }}

ExecStop=-{{ devture_systemd_docker_base_host_command_sh }} -c '{{ devture_systemd_docker_base_host_command_docker }} stop -t {{ devture_systemd_docker_base_container_stop_grace_time_seconds }} {{ promtail_identifier }} 2>/dev/null || true'
ExecStop=-{{ devture_systemd_docker_base_host_command_sh }} -c '{{ devture_systemd_docker_base_host_command_docker }} rm {{ promtail_identifier }} 2>/dev/null || true'
Restart=always
RestartSec=30
SyslogIdentifier={{ promtail_identifier }}

[Install]
WantedBy=multi-user.target
